# docs/FreeTAKHub/Video/mediamtx.service
[Unit]
After=network.target
[Service]
ExecStart=/usr/local/bin/mediamtx /usr/local/etc/mediamtx.yml
[Install]
WantedBy=multi-user.target